# systemd/reflicate.service
[Unit]
Description=Reflicate

[Service]
Type=oneshot
ExecStart=/bin/sh -c '~/.cargo/bin/reflicate ~'

[Install]
WantedBy=multi-user.target
